# minecraft-health-docker.service
[Unit]
Description=Minecraft Server Health Check Service (Docker)
After=network.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3.12 /home/ubuntu/minecraft-health-server-docker.py
Restart=always
RestartSec=10

# 環境変数の設定
Environment="HEALTH_CHECK_PORT=8080"
Environment="DOCKER_CONTAINER_NAME=paper"

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minecraft-health

[Install]
WantedBy=multi-user.target
